# rootfs/etc/services.d/crond/run
#!/usr/bin/env sh

exec /usr/sbin/cron -f -L 15
